# api-electoral.service
[Unit]
Description=API Electoral - FastAPI Application
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/var/www/html/apielectoral
Environment="PATH=/var/www/html/apielectoral/venv/bin"
Environment="ENVIRONMENT=production"
Environment="PORT=8000"
Environment="WORKERS=2"
EnvironmentFile=/var/www/html/apielectoral/.env

# Usar run.py en lugar de uvicorn directamente
ExecStart=/var/www/html/apielectoral/venv/bin/python run.py

Restart=always
RestartSec=5

# Limites de recursos
LimitNOFILE=65536

# Logs
StandardOutput=journal
StandardError=journal

# Variables de entorno para Playwright
Environment="PLAYWRIGHT_BROWSERS_PATH=/var/www/html/apielectoral/.playwright"

[Install]
WantedBy=multi-user.target
